# .backup-2017-08-05/.config/systemd/user/default.target
[Unit]
Description=Napalm's default.target
Wants=dbus.socket xinitrc.socket awesome.service

[Install]
Alias=default.target
